# docker/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true

; [program:laravel-worker]
; process_name=%(program_name)s_%(process_num)02d
; command=php /var/www/artisan queue:work --sleep=3 --tries=3 --max-time=3600
; autostart=true
; autorestart=true
; stopasgroup=true
; killasgroup=true
; user=laravel
; numprocs=2
; redirect_stderr=true
; stdout_logfile=/var/www/storage/logs/worker.log
; stopwaitsecs=3600

[program:cron]
command=cron -f
autostart=true
autorestart=true
